# bin/webssh2.service
[Unit]
Description=webssh2 Server daemon
[Service]
ExecStart=/usr/share/node-v10.6.0/bin/node /usr/share/quick/extend/webssh2/index.js
ExecStop=/usr/bin/pkill -f /usr/share/quick/extend/webssh2/index.js
Type=simple
KillMode=process
Restart=on-failure
RestartSec=42s

[Install]
WantedBy=multi-user.target
